# omnia/files/etc/hotplug.d/net/60-septentrio
[ "$ACTION" = "add" ] || exit 0
[[ "$INTERFACE" = "usb"* ]] || exit 0

idVendor="$(cat "/sys${DEVPATH}/../../../idVendor")"
[ "$idVendor" = "152a" ] || exit 0

/sbin/ip link set "$INTERFACE" name "usb_septentrio"
